# misc/nginx/app/logrotate.d/cron.daily/logrotate
#!/bin/sh

test -x /usr/sbin/logrotate || exit 0
/usr/sbin/logrotate -s /data/log/status /etc/logrotate.conf
